# logstash-salt/files/logstash.ini
[program:logstash]
command=/webdata/opt/local/logstash-2.0/bin/logstash agent -f /webdata/opt/local/logstash-2.0/nginx.conf
autostart=true
autorestart=true
startretries=3
stderr_logfile=/tmp/logstash.log
stdout_logfile=/tmp/logstash.log
